# system_files/usr/lib/systemd/user/polkit-agent.service
[Unit]
Description=Polkit authentication agent (niri only)
PartOf=graphical-session.target
After=graphical-session.target

[Service]
Type=simple
ExecCondition=/usr/bin/sh -c '[ "${XDG_CURRENT_DESKTOP:-}" = "niri" ] || [ "${XDG_CURRENT_DESKTOP:-}" = "sway" ]'
ExecStart=/usr/libexec/polkit-mate-authentication-agent-1
Restart=on-failure
RestartSec=2

[Install]
WantedBy=graphical-session.target
