# system_files/usr/lib/systemd/system/thinkfan.service.d/barbatos.conf
[Unit]
# Only run on ThinkPads, and wait for the profile switch to write the config
ConditionPathExists=/proc/acpi/ibm/fan
ConditionPathExists=/run/thinkfan-active.yaml

[Service]
# Use the runtime config written by barbatos-fanprofile-switch.sh
ExecStart=
ExecStart=/usr/bin/thinkfan -c /run/thinkfan-active.yaml $THINKFAN_ARGS
