# system_files/usr/lib/systemd/system/barbatos-fanprofile-switch.service
[Unit]
Description=Switch thinkfan profile based on AC/battery state
ConditionPathExists=/proc/acpi/ibm/fan

[Service]
Type=oneshot
ExecStart=/usr/local/bin/barbatos-fanprofile-switch.sh
